# com/pipool.service
[Unit]
Description=Start front-end & back-end node.js server for pipool/pool-control in /lib/systemd/system/
After=multi-user.target